# kai-services/streaming/docker/supervisord.conf
; kai 配信スタックのプロセス管理（起動順は priority で制御）
; 設計: docs/kai/design/streaming.md §5
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log

; 永続 X サーバー（:0, dummy ドライバ）
; -ac: コンテナ内限定の X ソケットのためアクセス制御を無効化（外部には公開されない）
[program:xorg]
command=/usr/lib/xorg/Xorg :0 -ac -noreset -nolisten tcp
user=kai
environment=HOME="/home/kai"
priority=10
autorestart=true

; PulseAudio（null-sink kai_speaker は /etc/pulse/default.pa.d/ で定義）
[program:pulseaudio]
command=/usr/bin/pulseaudio --daemonize=no --exit-idle-time=-1
user=kai
environment=HOME="/home/kai",XDG_RUNTIME_DIR="/run/user/1000"
priority=15
autorestart=true

; XFCE デスクトップセッション
[program:desktop]
command=/bin/bash -c 'until xdpyinfo >/dev/null 2>&1; do sleep 1; done; exec dbus-launch --exit-with-session startxfce4'
user=kai
environment=HOME="/home/kai",DISPLAY=":0",XDG_RUNTIME_DIR="/run/user/1000"
priority=20
autorestart=true

; リモート GUI（ポートは compose で host 127.0.0.1 にのみ publish）
[program:x11vnc]
command=/usr/bin/x11vnc -display :0 -forever -shared -noxdamage -rfbport 5900 -rfbauth /home/kai/.vnc/passwd
user=kai
environment=HOME="/home/kai"
priority=30
autorestart=true
