# kai-services/koe-bridge/koe-bridge.service
# koe-bridge（hermes auxiliary → OpenAI 互換中継）。systemd --user テンプレート。
# install.sh が @REPO_DIR@ を実パスに置換して ~/.config/systemd/user/ に配置する。
[Unit]
Description=kai koe-bridge (OpenAI-compatible relay to hermes auxiliary LLM)
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=@REPO_DIR@
Environment=KOE_BRIDGE_PORT=8930
# 既定は 127.0.0.1（loopback）。Tailscale 越しに使う場合は bind を広げ、
# 共有トークンを必ず設定する（トークン無しの非 loopback bind は起動拒否）。
# トークンは unit に直書きせず EnvironmentFile に置く（Issue #77 C1）:
#   ~/.config/kai/koe-bridge.env に KOE_BRIDGE_TOKEN=<ランダム値> を書く
# Environment=KOE_BRIDGE_BIND=0.0.0.0
# EnvironmentFile=%h/.config/kai/koe-bridge.env
ExecStart=@REPO_DIR@/.venv/bin/python3 @REPO_DIR@/kai-services/koe-bridge/koe_bridge.py
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
